# config/systemd/user/whatidid-org-export.service
[Unit]
Description=Exports Org project WhatIDid and pushes Git to GH Pages

[Service]
ExecStart=/home/dale/.dotfiles/bin/whatidid-org-export.sh
